# ansible/roles/workers/templates/scorinator_upstart.conf
description "scorinator service"
start on runlevel [2345]
stop on runlevel [06]

# ran based on ini file, die on term swaps SIGQUIT/SIGTERM behavior
#exec #need to put the worker binary path here
